# com.onyx.quadcopter/src/main/resources/onyx.service
[Unit]
Description=Onyx Quadcopter service

[Service]
ExecStart=/opt/onyx/Onyx/onyx_start
Restart=on-failure

[Install]
WantedBy=multi-user.target
