# com.onyx.quadcopter/src/main/resources/systemd/fan.service
[Unit]
Description=Fan Control

[Service]
ExecStart=/opt/onyx/fan.sh

[Install]
WantedBy=default.target
