# com.onyx.quadcopter/src/main/resources/systemd/onyx-fan.service
[Unit]
Description=Fan Control

[Service]
ExecStart=/opt/onyx/onyx-fan

[Install]
WantedBy=default.target
