# device-files/systemd/onyx-system.service
[Unit]
Description=Onyx Quadcopter service

[Service]
ExecStart=/opt/onyx/onyx-start
Restart=on-failure

[Install]
WantedBy=multi-user.target
